# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Planificateur Laravel (tâches programmées : sauvegardes, etc.).
# NB : en multi-réplicas, ne le lancer que dans une seule instance (RUN_SCHEDULER=false ailleurs)
# pour éviter les exécutions en double.
[program:scheduler]
command=php /var/www/html/artisan schedule:work
autostart=%(ENV_RUN_SCHEDULER)s
autorestart=true
startretries=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Worker de file d'attente (QUEUE_CONNECTION=database par défaut).
[program:queue]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --rest=1
autostart=%(ENV_RUN_QUEUE)s
autorestart=true
startretries=10
stopwaitsecs=3700
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
